# tools/scripts/com.mantissa.plot-bot-dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mantissa.plot-bot-dashboard</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/polansk/Developer/mantissa/code/plot-bot/tools/scripts/build_cycle_dashboard.py</string>
        <string>--serve</string>
        <string>--port</string>
        <string>3000</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/polansk/Developer/mantissa/code/plot-bot</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/polansk/Developer/mantissa/code/plot-bot/logs/dashboard_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/polansk/Developer/mantissa/code/plot-bot/logs/dashboard_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/polansk</string>
    </dict>
</dict>
</plist>
