# tools/scripts/com.mantissa.plot-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mantissa.plot-bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/polansk/Developer/mantissa/code/plot-bot/tools/scripts/runner.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/polansk/Developer/mantissa/code/plot-bot</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>60</integer>

    <key>StandardOutPath</key>
    <string>/Users/polansk/Developer/mantissa/code/plot-bot/logs/launchd_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/polansk/Developer/mantissa/code/plot-bot/logs/launchd_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/polansk</string>
        <key>PLOT_BOT_KEYCHAIN_PW</key>
        <string>SET_ON_INSTALL</string>
    </dict>
</dict>
</plist>
